# PAI/USER/FLOWS/com.erikschwartz.substack-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.erikschwartz.substack-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/bun</string>
    <string>/Users/erikschwartz/.claude/PAI/USER/FLOWS/F_SUBSTACK_MONITOR.ts</string>
    <string>--output</string>
    <string>/Users/erikschwartz/.claude/MEMORY/LOGS/substack-digest-latest.md</string>
  </array>

  <!-- Run every 6 hours -->
  <key>StartInterval</key>
  <integer>21600</integer>

  <key>StandardOutPath</key>
  <string>/Users/erikschwartz/.claude/MEMORY/LOGS/substack-monitor.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/erikschwartz/.claude/MEMORY/LOGS/substack-monitor.log</string>

  <key>ProcessType</key>
  <string>Background</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/erikschwartz</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/erikschwartz/.local/bin</string>
  </dict>
</dict>
</plist>
